# init/turnip.conf
#turnip
start on runlevel [2345]
stop on runlevel [016]

respawn

exec turnip /etc/turnip/settings.conf
